# .ddev/web-build/laravel-worker.conf
[program:laravel-worker]
process_name=%(program_name)s
command=php /var/www/html/artisan queue:work
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/queue.log
stopwaitsecs=3600